# misc/crone/renew_cert.cron
# Renew cert files

0 3  1  *  *  root certbot renew
0 3  1  *  *  root  nginx -s reload